# tc-executor.service
[Unit]
Description=TC Executor
After=network.target

[Service]
ExecStart=/path/to/your/script.sh

[Install]
WantedBy=multi-user.target
